# bluetooth_sensors/bluetooth-sensors.service
# copy this file to /etc/systemd/system/bluetooth-sensors.service

[Unit]
Description=Bluetooth Sensors Data Collection Service
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=kevin
WorkingDirectory=/home/kevin/github/GoHome
ExecStart=/home/kevin/github/GoHome/bluetooth_sensors/run.sh
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Optional: Adjust resource limits if needed
# LimitNOFILE=4096

RuntimeMaxSec=86400  # Optional: force restart after 24 hours if timer fails

[Install]
WantedBy=multi-user.target